# cron/crontab
1 0 * * * * php /app/handlers/producer.php > /proc/1/fd/1